# docker/images/sandbox-slim/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/home/dsb/supervisord.pid
childlogdir=/home/dsb

[program:chromium]
command=/usr/bin/chromium --no-sandbox --disable-setuid-sandbox --disable-gpu --disable-dev-shm-usage --remote-debugging-port=9222 --remote-debugging-address=0.0.0.0 --user-data-dir=/tmp/.chromium-cdp --headless --no-first-run --no-default-browser-check --silent-debugger --enable-logging --log-level=0
autorestart=true
priority=1
user=dsb
environment=CHROMIUM_PORT="9222",NODE_NO_WARNINGS="1"
stdout_logfile=/dev/null
stderr_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:tool_proxy]
command=/usr/local/bin/python3 /opt/tools/tool_proxy.py
autostart=true
autorestart=true
priority=50
user=root
directory=/opt/tools
stdout_logfile=/home/dsb/tool_proxy.log
stderr_logfile=/home/dsb/tool_proxy_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
