# docker/images/sandbox/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/home/dsb/supervisord.log
pidfile=/home/dsb/supervisord.pid
childlogdir=/home/dsb

[program:xvfb]
command=/bin/bash -c "exec /usr/bin/Xvfb :1 -screen 0 ${VNC_RESOLUTION:-2560x1440}x24 -ac -extension MIT-SHM"
autorestart=true
priority=10
environment=DISPLAY=":1"

[program:xfce]
command=/bin/bash -c "sleep 3 && exec /usr/bin/dbus-run-session /usr/bin/xfce4-session --display=:1"
autorestart=true
priority=15
user=dsb
environment=DISPLAY=":1",HOME="/home/dsb",USER="dsb"

[program:x11vnc]
command=/bin/bash -c "sleep 3 && exec /usr/bin/x11vnc -display :1 -forever -nopw -shared -rfbport 5901 -wait 5 -noshm"
autorestart=true
priority=20
user=dsb
environment=DISPLAY=:1

[program:chromium]
; Proxy: use wrapper so --proxy-server/--proxy-bypass-list match HTTP(S)_PROXY/NO_PROXY (Chromium + CDP + crawl4ai).
command=/opt/tools/chromium-launch.sh
autorestart=true
priority=1
user=dsb
environment=DISPLAY=:1

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc/ 6080 localhost:5901
priority=40
user=dsb
autorestart=true
environment=DISPLAY=:1

[program:tool_proxy]
command=/usr/local/bin/python3 /opt/tools/tool_proxy.py
autostart=true
autorestart=true
priority=1
user=root
directory=/opt/tools
stdout_logfile=/home/dsb/tool_proxy.log
stderr_logfile=/home/dsb/tool_proxy_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
